# conf/systemd/agora.service
[Unit]
Description=agora daemon
Requires=agora.socket
After=network.target

[Service]
User=agora
Group=www-data
WorkingDirectory=/home/agora/agora.gdn
ExecStart=/home/agora/agora.gdn/.venv/bin/gunicorn \
          --config /home/agora/agora.gdn/conf/gunicorn/gunicorn.conf.py \
          agora.wsgi

[Install]
WantedBy=multi-user.target