# conf/systemd/agora.socket
[Unit]
Description=agora socket

[Socket]
ListenStream=/run/agora.sock

[Install]
WantedBy=sockets.target
